# target/linux/tensor/base-files/etc/init.d/load-modules
#!/bin/sh /etc/rc.common
# Init script to load modules in order from /etc/modules.load

START=10
STOP=90

start() {
  [ -f /etc/modules.load ] || return 0
  echo "[modloader] loading kernel modules from /etc/modules.load"
  while read module; do
    mod="$(basename "$module")"
    echo "[modloader] insmod $mod"
    insmod "/lib/modules/$mod" 2>/dev/null || echo "Failed to load $mod"
  done < /etc/modules.load
}
